# autostart/Flameshot.desktop
[Desktop Entry]
Type=Application
Name=flameshot
Icon=flameshot
Exec=flameshot
Terminal=false
X-GNOME-Autostart-enabled=true
